# etc/supervisord.conf
[program:czcake8080]
command=python /root/app/cakeshop-master/manager.py
autostart=true
stderr_logfile = /var/log/nginx/app.log
stdout_logfile = /var/log/nginx/app.log
